# deploy/systemd/ccc-storage-mountd.service
[Unit]
Description=CCC layered storage mount daemon
Documentation=file:/opt/ccc-storage/deploy/README.md
After=network-online.target remote-fs.target
Wants=network-online.target

[Service]
Type=simple
Environment=CCC_STORAGE_MOUNTD_CONFIG=/etc/ccc-storage/mountd.toml
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
RuntimeDirectory=ccc-storage
RuntimeDirectoryMode=0755
ExecStart=/usr/bin/env ccc-storage mountd --config ${CCC_STORAGE_MOUNTD_CONFIG}
Restart=on-failure
RestartSec=2s
# Host-level daemon: needs local FUSE/mount capability on nodes where real FUSE
# adapters are enabled. Keep this opt-in via deployment, not container default.
DeviceAllow=/dev/fuse rw
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_SETPCAP CAP_CHOWN CAP_DAC_OVERRIDE CAP_FOWNER CAP_MKNOD CAP_SYS_CHROOT
AmbientCapabilities=CAP_SYS_ADMIN
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=full
ProtectHome=false

[Install]
WantedBy=multi-user.target
