# deploy/systemd/ccc-storage-mountd.service
[Unit]
Description=CCC layered storage mount daemon
Documentation=file:/opt/ccc-layered-storage/deploy/README.md
After=network-online.target remote-fs.target
Wants=network-online.target

[Service]
Type=simple
Environment=CCC_NFS_ROOT=/storage/.ccc-layered
Environment=CCC_NODE_RUN_DIR=/run/ccc-layered
Environment=CCC_MOUNTD_SOCK=/run/ccc-layered/mountd.sock
Environment=CCC_MANAGED_PARENT=/managed/datasets
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
RuntimeDirectory=ccc-layered
RuntimeDirectoryMode=0755
ExecStart=/usr/bin/env ccc-storage mountd --nfs-root ${CCC_NFS_ROOT} --run-dir ${CCC_NODE_RUN_DIR} --socket ${CCC_MOUNTD_SOCK} --managed-parent ${CCC_MANAGED_PARENT}
Restart=on-failure
RestartSec=2s
# Host-level daemon: needs local FUSE/mount capability on nodes where real FUSE
# adapters are enabled. Keep this opt-in via deployment, not container default.
DeviceAllow=/dev/fuse rw
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_SETPCAP CAP_CHOWN CAP_DAC_OVERRIDE CAP_FOWNER CAP_MKNOD CAP_SYS_CHROOT
AmbientCapabilities=CAP_SYS_ADMIN
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=full
ProtectHome=false

[Install]
WantedBy=multi-user.target
